# rootfs/etc/crontabs/selfoss
05 0 * * * /usr/sbin/logrotate /etc/logrotate.d/selfoss -s /selfoss/data/logs/logrotate.status